# service/xwayncd.service
[Unit]
Description=xwayncd

[Service]
ExecStart=/bin/bash /opt/xwaync/start_xwayncd.sh

[Install]
WantedBy=multi-user.target