# deploy/xwayncd.service
[Unit]
Description=xwayncd

[Service]
ExecStartPre=/bin/sleep 20
ExecStart=/opt/xwaync

[Install]
WantedBy=multi-user.target