# deploy/xwayncd.service
[Unit]
Description=xwayncd

[Service]
ExecStart=/opt/xwaync

[Install]
WantedBy=multi-user.target